# rootfs/etc/udev/rules.d/80-monitor-dock.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/monitor-hotplug.sh"
